# boincwatch/etc/init.d/boincwatch
#!/sbin/runscript

depend() {
	need boinc
}

start()
{
	ebegin "starting boincwatch"
	start-stop-daemon --start -p /var/run/boincwatch.pid -mb --chdir /var/lib/boinc/ \
			--exec /bin/bash /var/lib/boinc/boincwatch.sh --chuid boinc:boinc \
			--nicelevel 18 2>&1 
	eend $?
}

stop()
{
	ebegin "stopping boincwatch"
	start-stop-daemon --stop --pidfile /var/run/boincwatch.pid --signal INT
	eend $?
}
